# service/nc-backup.service
[Unit]
Description=Backup Nextcloud Configuration

[Service]
ExecStart=/opt/nextcloud/backup
